# crontab/crontab
* * * * * DISPLAY=:0.0 sh /home/indium/.battery_warning.sh
